# crontab/root
MIN */6 * * *              /root/mcss/mcss
30 04 * * *               /root/mcss/anti-malware
